# Code/mcvpn/rejoin-openstack.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec ~/devstack/rejoin-stack.sh

